# flask-backend.service
[Unit]
Description=Gunicorn instance to serve main
After=network.target